# deploy/agent-tree-rl.service
[Unit]
Description=Agent Tree RL decision controller
Documentation=file:/opt/agent-tree-rl/docs/RUNBOOK.md
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
DynamicUser=yes
WorkingDirectory=/opt/agent-tree-rl
StateDirectory=agent-tree-rl
StateDirectoryMode=0750
RuntimeDirectory=agent-tree-rl
RuntimeDirectoryMode=0750

Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=AGENT_TREE_RL_HOST=127.0.0.1
Environment=AGENT_TREE_RL_PORT=8081
Environment=AGENT_TREE_RL_DATA_DIR=/var/lib/agent-tree-rl
Environment=AGENT_TREE_RL_LOG_FORMAT=json
Environment=AGENT_TREE_RL_REQUIRE_AUTH=true
Environment=AGENT_TREE_RL_ALLOW_SAMPLE_BENCHMARK=false
Environment=AGENT_TREE_RL_REQUIRE_SEPARATION_OF_DUTIES=true
Environment=AGENT_TREE_RL_MAX_OPERATIONAL_WORKERS=2
Environment=AGENT_TREE_RL_SHUTDOWN_GRACE_SECONDS=30
Environment=AGENT_TREE_RL_SHUTDOWN_CANCEL_SECONDS=5
Environment=AGENT_TREE_RL_RECEIPT_KEYS_FILE=%d/receipt-keys.json
Environment=AGENT_TREE_RL_ADMIN_TOKEN_FILE=%d/api-tokens.json
Environment=AGENT_TREE_RL_BACKUP_KEYS_FILE=%d/backup-keys.json
Environment=AGENT_TREE_RL_BENCHMARK_FILE=%d/policy-benchmark.json
Environment=AGENT_TREE_RL_BENCHMARK_SIGNING_KEY_FILE=%d/benchmark-signing.key
Environment=AGENT_TREE_RL_ALLOWED_COMMANDS=/opt/agent-tree-rl/venv/bin/agent-tree-rl-evidence-probe
Environment=AGENT_TREE_RL_ALLOWED_CWD_ROOTS=/var/lib/agent-tree-rl
EnvironmentFile=-/etc/agent-tree-rl/service.env

# Prefer LoadCredentialEncrypted in a production image. LoadCredential keeps
# plaintext credentials out of the process environment and unit file.
LoadCredential=receipt-keys.json:/etc/agent-tree-rl/receipt-keys.json
LoadCredential=api-tokens.json:/etc/agent-tree-rl/api-tokens.json
LoadCredential=benchmark-signing.key:/etc/agent-tree-rl/benchmark-signing.key
LoadCredential=policy-benchmark.json:/etc/agent-tree-rl/policy-benchmark.json
LoadCredential=backup-keys.json:/etc/agent-tree-rl/backup-keys.json

ExecStart=/opt/agent-tree-rl/venv/bin/python -m agent_tree_rl.cli serve
Restart=on-failure
RestartSec=2s
TimeoutStartSec=30s
# Keep the supervisor deadline longer than the application's bounded drain.
TimeoutStopSec=40s
KillSignal=SIGTERM
KillMode=mixed

NoNewPrivileges=yes
CapabilityBoundingSet=
AmbientCapabilities=
UMask=0027
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=yes
ProtectClock=yes
ProtectHostname=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectProc=invisible
ProcSubset=pid
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
IPAddressDeny=any
IPAddressAllow=localhost
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
RemoveIPC=yes
ReadWritePaths=/var/lib/agent-tree-rl
SocketBindDeny=any
SocketBindAllow=ipv4:tcp:8081
SocketBindAllow=ipv6:tcp:8081

LimitNOFILE=65536
TasksMax=256
MemoryMax=1G
MemoryHigh=768M
CPUQuota=200%

[Install]
WantedBy=multi-user.target
